# deploy/portfolio.service
# /etc/systemd/system/portfolio.service
# Runs the Astro SSR runtime (Home page + /api/* endpoints).
# Resolves through the `current` symlink so a promoted release is picked up on reload.
[Unit]
Description=Portfolio (Astro SSR)
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/portfolio/current
ExecStart=/usr/bin/node /var/www/portfolio/current/dist/server/entry.mjs
Environment=HOST=127.0.0.1
Environment=PORT=4321
EnvironmentFile=/var/www/portfolio/shared/.env
Restart=always
RestartSec=2
User=www-data
Group=www-data
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
